# current/3-settings-lxqt/debian/whiteroom-lxqt-livesettings.service
[Unit]
Description=LXQt cheatcodes in live media
ConditionPathExists=/usr/share/whiteroom-lxqt-livesettings/
After=whiteroom-common-livesettings.service

[Service]
Type=oneshot
ExecStart=/usr/share/whiteroom-lxqt-livesettings/lxqt-cheatcodes start

[Install]
WantedBy=graphical.target
